# container_files/etc/services.d/nginx/run
#!/usr/bin/with-contenv /bin/sh
nginx -g 'daemon off;'