# container_files/etc/services.d/shibboleth/run
#!/usr/bin/with-contenv /bin/sh

echo "Starting shibboleth"
# shibd -F
shibd -fF  -c /etc/shibboleth/shibboleth2.xml